# systemd-user/user/goimapnotify.timer
[Unit]
Description=Remote imap mail notify service

[Timer]
OnBootSec=30
Unit=goimapnotify@gmail.service

[Install]
WantedBy=timers.target
