# systemd-user/user/notmuch.timer
[Unit]
Description=Mail indexing service

[Timer]
OnBootSec=60
OnUnitActiveSec=60m
Unit=notmuch.service

[Install]
WantedBy=timers.target
